# reasoning-tools.service
[Unit]
Description=Reasoning Tools MCP Server (Dual: SSE + Streamable HTTP)
After=network.target

[Service]
Type=simple
User=heefoo
WorkingDirectory=/home/heefoo/reasoning-tools
ExecStart=/home/heefoo/reasoning-tools/reasoning-tools -transport=dual -port=9847 -base-url http://localhost:9847 -http-path /mcp
# Alternative (SSE only):
# ExecStart=/home/heefoo/reasoning-tools/reasoning-tools -transport=sse -port=9847 -base-url http://localhost:9847
# Alternative (Streamable HTTP only):
# ExecStart=/home/heefoo/reasoning-tools/reasoning-tools -transport=streamable-http -port=9847 -http-path /mcp
Restart=on-failure
RestartSec=5

# Enable streaming output channels (all can be enabled simultaneously)
Environment=MCP_STDERR_STREAM=true
Environment=MCP_LOGGING_STREAM=true
Environment=MCP_PROGRESS_STREAM=true

# Rate limiting for concurrent LLM requests (default: 2, 0 = unlimited)
Environment=LLM_MAX_CONCURRENT=1

# Environment variables for LLM provider
Environment=LLM_PROVIDER=zai
Environment=LLM_MODEL=glm-4.7
Environment=ZAI_BASE_URL=https://api.z.ai/api/paas/v4
Environment=ZAI_API_KEY=e988037e10444f42b30e54a746d49797.WUFfyb2pEPFGDsjp

[Install]
WantedBy=multi-user.target
